# 4.1/crontab
*/1 * * * * /usr/src/redmine/cronjob.sh >> /usr/src/redmine/log/cron.log 2>&1
